# sphinx/2.3.2/rootfs/etc/services.d/searchd/run
#!/usr/bin/execlineb -P

s6-setuidgid sphinx
searchd --config /etc/sphinx/sphinx.conf --nodetach
