# nginx-php-fpm/1.13-7.1/rootfs/etc/services.d/nginx/run
#!/usr/bin/execlineb -P

foreground { s6-echo "[services.d/nginx]" }
foreground { s6-mkdir -p /run/nginx }

nginx -g "daemon off;"
